# deploy/supervisor/conf.d/live.conf
[program:live]
command=%(ENV_HOME)s/projects/roxabi-live/.venv/bin/roxabi-live
directory=%(ENV_HOME)s/projects/roxabi-live
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/projects/roxabi-live/.venv/bin:%(ENV_PATH)s"
priority=100
autostart=false
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/roxabi-live/logs/live.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/roxabi-live/logs/live_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
