# deploy/systemd/live.service
[Unit]
Description=Roxabi Live — ops cockpit (FastAPI + uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/roxabi-live
EnvironmentFile=%h/projects/roxabi-live/.env
ExecStart=%h/projects/roxabi-live/.venv/bin/roxabi-live
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
SyslogIdentifier=live

[Install]
WantedBy=default.target
